# fs_overlay/etc/services.d/02-rakejob/run
#!/usr/bin/execlineb -P
bundle exec rake subscriber:sync_file RAILS_ENV=development
