# fs_overlay/etc/services.d/01-sidekiq/run
#!/usr/bin/execlineb -P
bundle exec sidekiq -C ./config/sidekiq.yml -e development
