# fs_overlay/etc/services.d/00-unicorn/run
#!/usr/bin/execlineb -P
bundle exec unicorn -E production -c ./config/unicorn.rb
